# rootfs/etc/services.d/openvpn/run
#!/usr/bin/with-contenv sh
cd /etc/openvpn/ && /usr/sbin/openvpn --config mullvad_${VPN_COUNTRY}_all.conf
